# deploy/supervisor.conf
; Copy to /etc/supervisor/conf.d/fuelHunter.conf on the server

[program:fuelHunter-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/fuelHunter/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/fuelHunter
autostart=true
autorestart=true
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/fuelHunter/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
user=www-data
